# scripts/launchd/ai.traffic.openclaw-sync-traffic-models.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.traffic.openclaw-sync-traffic-models</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__REPO__/scripts/openclaw-sync-traffic-models.py</string>
  </array>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OPENCLAW_CONFIG</key>
    <string>__HOME__/.openclaw/openclaw.json</string>
    <key>OPENCLAW_SYNC_RESTART</key>
    <string>0</string>
  </dict>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/openclaw-traffic-sync.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/openclaw-traffic-sync.log</string>
  <key>ProcessType</key>
  <string>Background</string>
  <key>Nice</key>
  <integer>10</integer>
</dict>
</plist>
